# meta-rdk-mtk-bpir4/recipes-connectivity/bt-host-cpc-hci-bridge/files/bt-host-cpc-hci-bridge.service
[Unit]
Description=Bluetooth CPC HCI Bridge
After=cpcd.service
Requires=cpcd.service

[Service]
Type=simple
ExecStart=/usr/bin/bt_host_cpc_hci_bridge
Restart=always
RestartSec=5
WorkingDirectory=/tmp

[Install]
WantedBy=multi-user.target
