# meta-rdk-mtk-bpir4/recipes-connectivity/openthread/files/otbr-agent.path
[Unit]
Description=Triggers otbr-agent.service when cpcd is started successfully.
Documentation= man:systemd.path

[Path]
DirectoryNotEmpty=/dev/shm/cpcd/cpcd_0/

[Install]
WantedBy=multi-user.target
